# systemd/lightshow-scheduler.service
[Unit]
Description=LightShow Scheduler (Baldrick Poller)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/lightshow-visitor-tracker-pi
Environment="PYTHONUNBUFFERED=1"
ExecStart=%h/lightshow-visitor-tracker-pi/.venv/bin/python -m app.scheduler
Restart=on-failure
User=%i

[Install]
WantedBy=multi-user.target
